# nginx/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi  --ini /docker_api/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx